# magicmirror.timer
[Unit]

[Timer]
OnStartupSec=10s

[Install]
WantedBy=timers.target
